# ansible-st2/roles/arteria_node/templates/supervisor_arteria_siswrap.conf.j2
[program:siswrap-ws]
process_name = siswrap-ws-%(process_num)s
command = /opt/arteria/arteria-siswrap-env/bin/siswrap-ws
              --config {{ arteria_siswrap_config }} --debug
numprocs = 1
