# ansible-st2/roles/arteria_node/templates/supervisor_arteria_bcl2fastq.conf.j2
[program:bcl2fastq-ws]
process_name = bcl2fastq-ws-%(process_num)s
command = /opt/arteria/arteria-bcl2fastq-env/bin/bcl2fastq-ws
              --config {{ arteria_bcl2fastq_config }}
numprocs = 1
